# scripts/factory/factory.timer
# scripts/factory/factory.timer
# Software Factory — systemd USER timer. wakeup.sh loops internally while the
# queue is non-empty (idle-retick), so this timer fires only as a fallback.
# Enable via `task factory:autopilot:install`.
[Unit]
Description=Software Factory dispatcher tick (idle-retick loop; 5 min fallback timer)
Documentation=file:///home/patrick/Bachelorprojekt/scripts/factory/README.md

[Timer]
Unit=factory.service
# First tick 2 min after the user manager comes up.
OnBootSec=2min
# Fallback: re-arm 5 min after the previous tick service exits.
# In practice wakeup.sh loops while the queue is non-empty, so the timer only
# fires after the queue drains and the service exits cleanly.
OnUnitInactiveSec=5min
# Run a tick that was missed while the host was off, on next boot.
Persistent=true
# Smear ±30s (reduced since idle-retick shrinks the gap between busy ticks).
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
